# scripts/deploy/equilytic.service
# Equilytic — FastAPI 服务 + 飞书 Stream 机器人
# 部署：拷到 /etc/systemd/system/equilytic.service，按需改 User / 路径 / 端口
#
# 说明：
# - .env 由程序自身 load_dotenv() 读取（放在 WorkingDirectory 即可），不要用 EnvironmentFile，
#   雪球 cookie 等含 ; = 的值会让 systemd 解析失败。
# - --serve-only 会同时拉起飞书 Stream 机器人（FEISHU_STREAM_ENABLED=true 且凭证完整时）。
# - 仅用本地 Claude 时才需要 equilytic-shim.service；用远程 LLM key 可去掉 Requires/After 中的 equilytic-shim。
[Unit]
Description=Equilytic (FastAPI API + Feishu bot)
After=network-online.target equilytic-shim.service
Wants=network-online.target
Requires=equilytic-shim.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/equilytic
ExecStart=/opt/equilytic/.venv/bin/python main.py --serve-only --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
